# cSploit/jni/firebird/builds/install/arch-specific/linux/firebird-classic.socket.in
[Unit]
Description=Firebird Classic Activation Socket
Conflicts=firebird-superserver.service

[Socket]
ListenStream=3050
Accept=true

[Install]
WantedBy=sockets.target
