# cSploit/jni/firebird/builds/install/arch-specific/linux/firebird-superserver.service.in
[Unit]
Description=Firebird Database Server ( SuperServer )
After=syslog.target network.target
Conflicts=firebird-classic.socket

[Service]
User=firebird
Group=firebird
Type=forking
PIDFile=/var/run/firebird/default.pid
ExecStart=@FB_SBINDIR@/fbguard -pidfile /var/run/firebird/default.pid -daemon -forever
StandardError=syslog

[Install]
WantedBy=multi-user.target
