# cSploit/jni/firebird/builds/install/arch-specific/linux/firebird-classic.service.in
[Unit]
Description=Firebird Database Server ( Classic )
After=local-fs.target

[Service]
User=firebird
Group=firebird
ExecStart=@FB_SBINDIR@/firebird
StandardInput=socket
StandardError=syslog
